# systemd/mithril-ops-console.service
[Unit]
Description=Mithril-OS Ops Console (host)
After=network.target

[Service]
Type=simple
User=mini-home-lab
WorkingDirectory=/mithril-os/apps/ops-console
EnvironmentFile=-/mithril-os/apps/ops-console/.env
ExecStart=/usr/bin/env node src/server.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
