# systemd/openclaw-usage-collector.service
[Unit]
Description=Collect OpenClaw usage telemetry snapshot
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/env bash /mithril-os/scripts/collect-openclaw-usage.sh
StandardOutput=journal
StandardError=journal
